# Arkserver/ark-backup-daily.service
# /etc/systemd/system/ark-backup-daily.service
[Unit]
Description=ARK Daily Backup (23:59)
Wants=ark-backup-daily.timer

[Service]
Type=oneshot
User=arkserver
ExecStart=/home/arkserver/arkserver/core/ark-backup.sh daily